# wot-gateway/init.d/victron-reader
### BEGIN INIT INFO
# Provides:          victron-reader
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start the Victron reader program at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin

. /lib/lsb/init-functions

set -e

DAEMON=/home/pi/development/LaMaisonVivante/electricity/victron-reader/target/release/victron-reader
SSD_OPTIONS="--exec $DAEMON -- --into-thing --thing-port 8081"

case $1 in
  start)
    log_daemon_msg "Starting $DESC"

    start-stop-daemon --start --background $SSD_OPTIONS

    log_progress_msg "${DAEMON##*/}"
    log_end_msg 0
  ;;

  stop)
    log_daemon_msg "Stopping $DESC"

    start-stop-daemon --stop $SSD_OPTIONS

    log_progress_msg "${DAEMON##*/}"
    log_end_msg 0
  ;;

  restart|force-reload)
    $0 stop
    sleep 1
    $0 start
  ;;

  status)
    status_of_proc "$DAEMON" "$DESC" && exit 0 || exit $?
  ;;

  *)
    echo "Usage: /etc/init.d/victron-reader-d {start|stop|restart|force-reload|status}" >&2
    exit 1
  ;;
esac

exit 0